# scripts/valheim-backup.timer
[Unit]
Description=Valheim backup timer

[Timer]
OnCalendar=*:0/30
Unit=valheim-backup.service

[Install]
WantedBy=timers.target
